# webserver/webservices.service
[Unit]
Description=LaMachine Webservices

[Service]
Type=forking
ExecStart=/usr/src/LaMachine/startwebservices.sh
ExecStop=killall -s 9 uwsgi
User=vagrant
Group=vagrant

[Install]
WantedBy=multi-user.target
